# storm-master.conf
[supervisord]
nodaemon=true

[program:nimbus]
command=/opt/storm/bin/storm nimbus
redirect_stderr=true

[program:stormui]
command=/opt/storm/bin/storm ui
redirect_stderr=true

