# storm-worker.conf
[supervisord]
nodaemon=true

[program:stormworker]
command=/opt/storm/bin/storm supervisor
redirect_stderr=true
